# port-forwarding/conf.d/http.conf
[program:http]
command=socat TCP-LISTEN:80,reuseaddr,fork TCP:%(ENV_HOST)s:39880
