# port-forwarding/conf.d/https.conf
[program:https]
command=socat TCP-LISTEN:443,reuseaddr,fork TCP:%(ENV_HOST)s:39843
