# port-forwarding/conf.d/minecraft.conf
[program:minecraft]
command=socat TCP-LISTEN:25565,reuseaddr,fork TCP:%(ENV_HOST)s:39865
