# port-forwarding/conf.d/mysql.conf
[program:mysql]
command=socat TCP-LISTEN:39806,reuseaddr,fork TCP:%(ENV_HOST)s:39806
