# port-forwarding/conf.d/ssh-gitlab.conf
[program:ssh-gitlab]
command=socat TCP-LISTEN:2289,reuseaddr,fork TCP:%(ENV_HOST)s:39889

