# port-forwarding/conf.d/vnc.conf
[program:vnc]
command=socat TCP-LISTEN:5900,reuseaddr,fork TCP:%(ENV_HOST)s:39859
